# services/plasmalogin.sh.in
#!/bin/sh
# non-systemd script for plasmalogin
# Provides: plasmalogin
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# XTODO: add proper dependency for getty
# Keybits: NO

. /etc/rc.subr

name="plasmalogin"
rcvar=plasmalogin_enable
command="@CMAKE_INSTALL_FULL_BINDIR@/plasmalogin"
pidfile="/run/plasmalogin/plasmalogin.pid"
sig_stop=TERM
sig_reload=USR1

load_rc_config $name

# Set environment variables for the daemon
export PLASMALOGIN_GREETER_EXEC=@CMAKE_INSTALL_FULL_BINDIR@/startplasma-login-x11

# Create required directories and log file before starting
plasmalogin_start() {
    mkdir -p /var/log/sonic
    chown root:root /var/log/sonic
    chmod 777 /var/log/sonic
    # Rotate log file
    [ -f /var/log/sonic/loginmanager.log.old ] && rm -f /var/log/sonic/loginmanager.log.old
    [ -f /var/log/sonic/loginmanager.log ] && mv /var/log/sonic/loginmanager.log /var/log/sonic/loginmanager.log.old
    touch /var/log/sonic/loginmanager.log
    chown plasmalogin:plasmalogin /var/log/sonic/loginmanager.log
    chmod 666 /var/log/sonic/loginmanager.log
    /bin/mkdir -p /run/plasmalogin
    chown root:plasmalogin /run/plasmalogin
    chmod 711 /run/plasmalogin
}

run_rc_command "$1"
