# services/plasmalogin.sh.bsd.in
#!/bin/sh
# BSD script for plasmalogin
# PROVIDE: plasmalogin
# REQUIRE: LOGIN dbus syslogd
# BEFORE: getty
# KEYWORD: shutdown

. /etc/rc.subr

name="plasmalogin"
rcvar="plasmalogin_enable"

load_rc_config $name

# Binary and state directory paths
command="@CMAKE_INSTALL_FULL_BINDIR@/plasmalogin"

# Signals
sig_stop="TERM"
sig_reload="USR1"

run_rc_command "$1"
