# src/frontend/startkde/plasma-login-kwin_wayland.service.in
[Unit]
Description=KDE Window Manager (Login Manager Version)
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kwin_wayland --no-lockscreen --no-global-shortcuts --no-kactivities --inputmethod plasma-keyboard --locale1
SuccessExitStatus=15
BusName=org.kde.KWin
Slice=session.slice
