# services/sysvinit/plasmalogin.in
#!/bin/sh
# sysvinit script for plasmalogin
# Provides: plasmalogin
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# XTODO: add proper dependency for getty
# Keybits: NO

. /etc/rc.subr

name="plasmalogin"
rcvar=plasmalogin_enable
command="@CMAKE_INSTALL_FULL_BINDIR@/plasmalogin"
pidfile="/run/plasmalogin/plasmalogin.pid"
sig_stop=TERM
sig_reload=USR1

load_rc_config $name

# Set environment variables for the daemon
export PLASMALOGIN_GREETER_EXEC=@CMAKE_INSTALL_FULL_BINDIR@/startplasma-login-x11

run_rc_command "$1"
